# ha-addon-ait-heatpump/rootfs/etc/s6-overlay/s6-rc.d/ha-addon/run
#!/command/with-contenv bashio

exec /app//HA.Addon.Ait.Heatpump --urls "http://*:5000"