# ha-addon-azure-backup/rootfs/etc/s6-overlay/s6-rc.d/ha-addon/run
#!/command/with-contenv bashio

exec /app/HA.Addon.Azure.Backup --urls "http://*:5000"